# immortalwrt/build.sh
#!/bin/bash
set -euxo pipefail

echo "ðŸš€ å¼€å§‹å‡†å¤‡æž„å»ºçŽ¯å¢ƒ..."

# æ›´æ–° feeds
./scripts/feeds update -a
./scripts/feeds install -a

# æ¸…ç†æ—§ç¼“å­˜
make clean || true
rm -rf tmp/ || true

# ===============================
# è‡ªå®šä¹‰å®‰è£…åŒ…ï¼ˆæ—  PPP ç»„ä»¶ï¼‰
# ===============================
PACKAGES=""
PACKAGES="$PACKAGES curl wget ca-certificates"
PACKAGES="$PACKAGES luci luci-base luci-compat luci-app-firewall"
PACKAGES="$PACKAGES luci-i18n-base-zh-cn luci-i18n-firewall-zh-cn luci-i18n-package-manager-zh-cn"
PACKAGES="$PACKAGES luci-app-docker luci-app-ttyd luci-app-filebrowser"
PACKAGES="$PACKAGES luci-app-opkg openssh-sftp-server"
PACKAGES="$PACKAGES kmod-usb-storage block-mount e2fsprogs fdisk"

# ===============================
# å†™å…¥è‡ªåŠ¨æ‰©å®¹è„šæœ¬
# ===============================
mkdir -p files/etc/init.d
cat > files/etc/init.d/expand_rootfs <<'EOF'
#!/bin/sh /etc/rc.common
START=99
DESCRIPTION="Auto expand root filesystem on first boot"

start() {
    if [ ! -f /etc/expand_done ]; then
        echo "ðŸ”§ è‡ªåŠ¨æ‰©å±• eMMC åˆ†åŒº..."
        parted /dev/mmcblk1 resizepart 2 100%
        losetup /dev/loop0 /dev/mmcblk1p2
        e2fsck -f -y /dev/loop0
        resize2fs -f /dev/loop0
        sync
        echo "âœ… æ‰©å±•å®Œæˆï¼Œé‡å¯ç”Ÿæ•ˆ..."
        touch /etc/expand_done
        reboot
    fi
}
EOF
chmod +x files/etc/init.d/expand_rootfs

# ===============================
# æž„å»ºé•œåƒï¼ˆè°ƒå¤§åˆ†åŒºç©ºé—´ï¼‰
# ===============================
echo "ðŸ§± å¼€å§‹æž„å»ºé•œåƒ..."
make image \
  PACKAGES="$PACKAGES" \
  FILES="files" \
  ROOTFS_PARTSIZE="1024" \
  V=s

# ===============================
# åŽ‹ç¼©å¹¶å‘å¸ƒ
# ===============================
OUTPUT_IMG=$(find bin/targets/ -name "*emmc-burn.img" | head -n 1)
if [ -f "$OUTPUT_IMG" ]; then
  echo "ðŸ“¦ åŽ‹ç¼©çº¿åˆ·åŒ…..."
  xz -T0 -z -9 "$OUTPUT_IMG"
fi

# ===============================
# ç”Ÿæˆæ›´æ–°è¯´æ˜Ž
# ===============================
mkdir -p ../release_note
cat > ../release_note/update.txt <<EOF
ðŸ†• æœ¬æ¬¡æ›´æ–°å†…å®¹ï¼š
- ç§»é™¤ PPPoE ç›¸å…³æ¨¡å—ï¼ˆppp-mod-pppoeã€kmod-pppoeã€pppï¼‰
- é€‚é…æ—è·¯ç”±æ¨¡å¼ï¼ˆDHCP è‡ªåŠ¨èŽ·å–ä¸Šçº§ IPï¼‰
- æ–°å¢žæ’ä»¶ï¼š
  - luci-app-docker
  - luci-app-ttyd
  - luci-app-filebrowser
- è‡ªåŠ¨æ‰©å±• eMMC å‰©ä½™ç©ºé—´
- é»˜è®¤å¯†ç ä¸ºç©º
EOF

echo "âœ… æž„å»ºå®Œæˆï¼"
echo "ðŸ“ çº¿åˆ·åŒ…æ–‡ä»¶: onecloud-immortalwrt-ext4-emmc-burn.img.xz"
echo "ðŸ“ æ›´æ–°è¯´æ˜Ž: release_note/update.txt"
