# immortalwrt/build.sh
#!/bin/bash
set -euxo pipefail

# =============================
# ImmortalWrt è‡ªå®šä¹‰å¿«é€Ÿæž„å»ºè„šæœ¬
# æ”¯æŒ OneCloud è‡ªåŠ¨æ‰©å±• overlay åˆ†åŒº
# é»˜è®¤ root æ— å¯†ç ï¼ŒIP 192.168.2.2
# =============================

# è‡ªå®šä¹‰è¦å®‰è£…çš„åŒ…
PACKAGES=""
PACKAGES="$PACKAGES curl"
PACKAGES="$PACKAGES luci-i18n-base-zh-cn"
PACKAGES="$PACKAGES luci-i18n-firewall-zh-cn"
PACKAGES="$PACKAGES luci-i18n-package-manager-zh-cn"
PACKAGES="$PACKAGES luci luci-app-opkg luci-app-docker luci-app-ttyd luci-app-filebrowser "

# æ ¹åˆ†åŒºå¤§å°ï¼ˆå•ä½ MBï¼‰
ROOTFS_PARTSIZE="512"

# åˆ›å»ºè‡ªå®šä¹‰æ–‡ä»¶ç»“æž„
mkdir -p files/etc/init.d/
mkdir -p files/etc/config/

# === è®¾ç½® LAN é»˜è®¤ IP ä¸º 192.168.2.2 ===
cat > files/etc/config/network <<'EOF'
config interface 'lan'
	option ifname 'eth0'
	option proto 'static'
	option ipaddr '192.168.2.2'
	option netmask '255.255.255.0'
	option gateway '192.168.2.1'
	option dns '223.5.5.5'
EOF

# === è®¾ç½® root å¯†ç ä¸ºç©º ===
mkdir -p files/etc
cat > files/etc/shadow <<'EOF'
root::0:0:99999:7:::
EOF

# === æ·»åŠ è‡ªåŠ¨æ‰©å±• overlay åˆ†åŒºè„šæœ¬ ===
cat > files/etc/init.d/expand-overlay <<'EOF'
#!/bin/sh /etc/rc.common
# è‡ªåŠ¨æ‰©å±• overlay åˆ†åŒºï¼Œä»…é¦–æ¬¡å¯åŠ¨æ‰§è¡Œ
START=99
STOP=10

start() {
    if [ -f /etc/expanded_done ]; then
        exit 0
    fi

    echo "å¼€å§‹æ‰©å±• overlay åˆ†åŒº..."
    parted /dev/mmcblk1 resizepart 2 100%
    losetup /dev/loop0 /dev/mmcblk1p2
    e2fsck -f -y /dev/loop0
    resize2fs -f /dev/loop0
    sync

    # æ‰©å±• Docker å­˜å‚¨ç›®å½•
    [ -d /opt/docker ] || mkdir -p /opt/docker
    mount -o bind /overlay/docker /opt/docker

    touch /etc/expanded_done
    echo "æ‰©å±•å®Œæˆï¼Œä¸‹æ¬¡å¯åŠ¨ä¸å†æ‰§è¡Œã€‚"
    reboot
}
EOF

chmod +x files/etc/init.d/expand-overlay

# === æ‰§è¡Œæž„å»º ===
echo "ðŸš€ å¼€å§‹æž„å»ºé•œåƒ..."
make image PACKAGES="$PACKAGES" FILES="files" ROOTFS_PARTSIZE="$ROOTFS_PARTSIZE"

echo "ðŸŽ‰ æž„å»ºå®Œæˆï¼é•œåƒå·²åŒ…å«ä»¥ä¸‹è‡ªå®šä¹‰åŠŸèƒ½ï¼š"
echo " - é»˜è®¤ IP: 192.168.2.2"
echo " - root å¯†ç ä¸ºç©ºï¼ˆç›´æŽ¥ç™»å½•ï¼‰"
echo " - è‡ªåŠ¨æ‰©å±• overlay åˆ†åŒºï¼ˆé¦–æ¬¡å¯åŠ¨è‡ªåŠ¨æ‰§è¡Œï¼‰"
echo " - é¢„è£…æ’ä»¶: luciã€dockerã€openclashã€ttydã€filebrowserã€nikki ç­‰"

